# centos8/etc/profile.d/php.sh
alias php='docker run --rm -it -v ${PWD}:/var/www/html -v laravel7trainning_mysql:/var/run/mysqld laravel7trainning_php74 php'
