# centos8/etc/profile.d/php.sh
alias php='docker run --rm -it -v ${PWD}:/var/www/html php:7.4-fpm php'
